# deployment/systemd/csky-api.service
[Unit]
Description=CartoSky API
After=network.target

[Service]
Type=simple
User=cartosky
Group=cartosky
UMask=0027
WorkingDirectory=/opt/cartosky
Environment=CARTOSKY_API_WORKERS=4
EnvironmentFile=-/etc/cartosky/api.env
Environment=PATH=/opt/cartosky/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/bin/sh -c 'exec /opt/cartosky/.venv/bin/uvicorn backend.app.main:app --host 127.0.0.1 --port 8200 --workers "${CARTOSKY_API_WORKERS:-${CARTOSKY_WORKERS:-4}}"'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
